# deploy/systemd/vaultwarden-secrets.service
[Unit]
Description=Vaultwarden Secrets Server
After=network.target

[Service]
Type=simple
User=vw-secrets
WorkingDirectory=/opt/vaultwarden-secrets
Environment="SECURITY_PROFILE=im-aware"
Environment="API_TOKEN_LXC200=..."
Environment="API_TOKEN_LXC202=..."
ExecStart=/usr/local/bin/bun run server/main.ts
Restart=on-failure

[Install]
WantedBy=multi-user.target
