# deploy/systemd/vw-backup.service
[Unit]
Description=Encrypted Vaultwarden Secrets state backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=mcp2cli
Group=ai-services
WorkingDirectory=/opt/vaultwarden-secrets
Environment="HOME=/var/lib/vaultwarden-secrets"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/etc/vaultwarden-secrets/backup.env
ExecStart=/usr/bin/env bun run scripts/backup.ts
NoNewPrivileges=true
PrivateTmp=true
