# deploy/systemd/vw-cred-proxy.service
[Unit]
Description=Vaultwarden Credential Proxy
After=network.target

[Service]
Type=simple
# Dedicated non-root service identity (issue #14).
User=vwsecrets
Group=vwsecrets
WorkingDirectory=/opt/vaultwarden-secrets

StateDirectory=vaultwarden-secrets
StateDirectoryMode=0700

# Explicit PATH/HOME — systemd minimal PATH breaks bw at /usr/local/bin/bw.
Environment=HOME=/var/lib/vaultwarden-secrets
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=VW_STATE_DIR=/var/lib/vaultwarden-secrets/state
Environment=BW_SESSION_FILE=/var/lib/vaultwarden-secrets/.bw-session
Environment=PROXY_PORT=3003
EnvironmentFile=-/opt/vaultwarden-secrets/.env.cred-proxy

ExecStart=/usr/local/bin/bun run server/cred-proxy.ts
Restart=on-failure
RestartSec=5

# Hardening (issue #14)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/vaultwarden-secrets

[Install]
WantedBy=multi-user.target
