# deploy/systemd/vw-deploy.timer
[Unit]
Description=Poll for new commits and deploy

[Timer]
OnCalendar=*:0/3
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
