# deploy/systemd/vw-snapshot.timer
[Unit]
Description=Periodic vault snapshot

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
