# deploy/systemd/vw-snapshot.service
[Unit]
Description=Vaultwarden Vault Snapshot
After=network.target

[Service]
Type=oneshot
User=vw-secrets
WorkingDirectory=/opt/vaultwarden-secrets
Environment="SECURITY_PROFILE=im-aware"
ExecStart=/usr/local/bin/bun run cli.ts snapshot
