# deploy/systemd/vw-snapshot.service
[Unit]
Description=Vaultwarden Vault Snapshot
After=network.target

[Service]
Type=oneshot
# Dedicated non-root service identity (issue #14).
User=mcp2cli
Group=ai-services
WorkingDirectory=/opt/vaultwarden-secrets

StateDirectory=vaultwarden-secrets
StateDirectoryMode=0700

# Explicit PATH/HOME — systemd minimal PATH breaks bw at /usr/local/bin/bw.
Environment=HOME=/var/lib/vaultwarden-secrets
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# State paths pinned into StateDirectory (DEP-3/DEP-6). The snapshot MUST write
# to the same config dir the servers read, and use the service-owned session +
# master key — not root's HOME-derived files.
Environment=VW_STATE_DIR=/var/lib/vaultwarden-secrets/state
Environment=VAULTWARDEN_SECRETS_DIR=/var/lib/vaultwarden-secrets/config
Environment=BW_SESSION_FILE=/var/lib/vaultwarden-secrets/.bw-session
Environment=MASTER_KEY_FILE=/var/lib/vaultwarden-secrets/.master-key
Environment=AUDIT_LOG_FILE=/var/lib/vaultwarden-secrets/audit.log
Environment=SECURITY_PROFILE=im-aware
EnvironmentFile=-/opt/vaultwarden-secrets/.env

ExecStart=/usr/local/bin/bun run cli.ts snapshot

# Hardening (issue #14)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/vaultwarden-secrets
